# systemd_scripts/rsync-backup.timer
[Unit]
Description=Esegui backup ogni 5 giorni

[Timer]
OnBootSec=10min
OnUnitActiveSec=5d
Persistent=true

[Install]
WantedBy=timers.target

